# crontab
0 1 * * * root sh /etc/bin/logcleaning.sh

