# crontab
0 0 * * * root sh /etc/bin/logcleaning.sh
